# clients/java/kurento-repository/kurento-repository-server/package/support-files/kurento-repo.service
[Unit]
Description=Kurento repository daemon
After=network.target

[Service]
User=nobody
ExecStart=/usr/bin/kurento-repo

[Install]
WantedBy=default.target