# clients/java/kurento-jsonrpc/kurento-jsonrpc-demo-server/package/support-files/kjrserver.service
[Unit]
Description=KJRServer daemon
After=network.target

[Service]
User=nobody
ExecStart=/usr/bin/kjrserver

[Install]
WantedBy=default.target